# config/rclone-daily.cron
# rclone daily backup
0 3 * * * root /usr/local/bin/rclone_daily.sh
